# riptide.service
[Unit]
Description=Riptide Crypto Node Manager
Documentation=https://github.com/deeep-network/riptide
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/app
ExecStart=/usr/bin/node /app/src/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=riptide

# Environment variables
Environment=NODE_ENV=production
Environment=RIPTIDE_SERVICE_NAME=riptide
Environment=CRYPTO_NODE_SERVICE=myria-node

# Security settings
NoNewPrivileges=false
PrivateTmp=false

# Resource limits
LimitNOFILE=65536
LimitNPROC=32768

[Install]
WantedBy=multi-user.target
