# systemd/myria.service
[Unit]
Description=Myria Node with Riptide SDK
Documentation=https://www.npmjs.com/package/@deeep-network/riptide
After=network.target
Wants=network.target

[Service]
Type=simple
User=myria
Group=myria
WorkingDirectory=/opt/myria
ExecStart=/usr/bin/node index.js
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=myria-riptide

# Environment variables
Environment=NODE_ENV=production
Environment=MYRIA_DATA_DIR=/var/lib/myria
Environment=MYRIA_LOG_DIR=/var/log/myria

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/myria /var/log/myria /opt/myria

[Install]
WantedBy=multi-user.target
